# centos7-nginx1.12-mysql5.6-php5.4/etc/supervisor/supervisord-mysql5.7-php5.4.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"

[program:php-fpm]
command=php-fpm --nodaemonize

[program:mysql]
command=/usr/local/bin/docker-entrypoint.sh mysqld