# superviosr-springboot-redis-mysql/supervisord-springboot-webapp.conf
[program:webapp]
command=java -Djava.security.egd=file:/dev/./urandom -jar /springboot/webapp.war
autostart=true
autorestart=true
user=root
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
